# systemd/minecraft.service
[Unit]
Description=Minecraft Paper Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=minecraft
Group=minecraft
WorkingDirectory=/opt/minecraft/server

ExecStart=/opt/minecraft/scripts/start.sh
ExecStop=/opt/minecraft/scripts/stop.sh

Restart=on-failure
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true
PrivateTmp=true
CapabilityBoundingSet=
AmbientCapabilities=

LimitNOFILE=65536
LimitNPROC=8192
LimitAS=infinity

PIDFile=/opt/minecraft/server/.pid

[Install]
WantedBy=multi-user.target
